# crontab
*/1 * * * * python3 /app/pulsar_data_collector.py >> /app/pulsar_data_collector.log 2>&1
