# crontab
*/1 * * * * /usr/bin/python3 pulsar_data_collector.py
